# gunicorn.service
[Unit]
Description=demo fastapi application
After=network.target

[Service]
User=xen
Group=xen

[install]
