# gunicorn.service
[Unit]
Description=demo fastapi application
After=network.target

[Service]


[install]
